# .config/systemd/user/hyprlock-update.service
[Unit]
Description=updates hyprlock information

[Service]
Type=oneshot
ExecStart=pkill -USR2 hyprlock
