# .config/systemd/user/hyprlock-update.service
[Unit]
Description=updates hyprlock information

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -USR2 hyprlock
