# .config/systemd/user/lock.target
[Unit]
Description=lock the current session
#Requires=graphical-session.target
Conflicts=unlock.target
