# .config/systemd/user/playerctl-pause-on-lock.service
[Unit]
Description=to pause all players
Before=sleep.target

[Service]
Type=oneshot
ExecStart=playerctl -sa pause

[Install]
WantedBy=lock.target
