# .config/systemd/user/screen-off.target
[Unit]
Description=turn off the screen
#Requires=graphical-session.target
Conflicts=wakeup.target
