# .config/systemd/user/sleep.target
[Unit]
Description=User-level sleep
Wants=lock.target
After=lock.target
