# .config/systemd/user/unlock.target
[Unit]
Description=unlock the current session
#Requires=graphical-session.target
Conflicts=lock.target
