# .config/systemd/user/wakeup.target
[Unit]
Description=stop sleep.target
#Requires=graphical-session.target
Conflicts=sleep.target
Conflicts=screen-off.target
