# .config/systemd/user/weather-update.timer
[Unit]
Description=timer to update cached weather

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
