# .config/systemd/user/weather-update.service
[Unit]
Description=update cached weather

[Service]
Type=oneshot
ExecStart=sh -c '/home/$USER/.local/bin/weather > /home/$USER/.cache/weather'
