# .config/systemd/user/weather-update.path
[Unit]
Description=monitor cached weather changes

[Path]
PathModified=/home/$USER/.cache/weather
Unit=hyprlock-update.service

[Install]
WantedBy=lock.target
