# dot_config/systemd/user/weather-update.timer
[Unit]
Description=timer to update cached weather

[Timer]
OnCalendar=*:0/15
RemainAfterElapse=false

[Install]
WantedBy=timers.target
