# .config/systemd/user/ssh-agent.service.d/override.conf
[Service]
Environment=SSH_ASKPASS=/usr/bin/ksshaskpass
